# s6-rc.d/fortivpn/run
#!/command/execlineb -P
with-contenv
/usr/bin/go-fortivpn-daemon
