# s6-rc.d/fortivpn/run
#!/bin/execlineb -P
expect /fortirun.expect