# s6-rc.d/fortivpn/run
#!/command/execlineb -P
/usr/bin/expect /fortirun.expect