# s6-rc.d/gost/run
#!/bin/execlineb -P
/usr/bin/gost -L "socks5://0.0.0.0:10800"
