# s6-rc.d/gost/run
#!/command/with-contenv bash
/usr/bin/gost -L "socks5+tls://${GOST_USERNAME}:${GOST_PASSWD}@0.0.0.0:10800"
